# dotfiles/.config/systemd/user/gtavibot.service
[Unit]
Description=Trailer 2 awaiting-bot

[Service]
ExecStart=python /home/work/bot/bot.py
Restart=on-failure

[Install]
WantedBy=default.target
