# dotfiles/.config/systemd/user/stagingsentinel.service
[Unit]
Description=Run my script every 2 minutes

[Service]
ExecStart=/home/work/staging_sentinel.sh

